# vm_deployment/noc-configmaker.service
[Unit]
Description=NOC Config Maker Backend Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__INSTALL_DIR__
Environment="PATH=__INSTALL_DIR__/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=__INSTALL_DIR__/venv/bin/python3 __INSTALL_DIR__/api_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__INSTALL_DIR__/secure_data

[Install]
WantedBy=multi-user.target
